# watchdog/watchdog-pinger.service
[Service]
Type=simple
ExecStart=/usr/local/bin/watchdog-pinger.sh
PIDFile=/run/watchdog-pinger.pid

[Install]
WantedBy=default.target
